# supervisor/db_mon.ini
[program:db_mon]
command=/home/jh/code/.venv/bin/python /home/jh/code/rec-crypto-data/rcd/db_mon.py
autostart=true
autorestart=true
stderr_logfile=/var/log/db_mon.err.log
stdout_logfile=/var/log/db_mon.out.log