# supervisor/fluent-bit.ini
[program:fluent-bit]
command=/usr/bin/fluent-bit -c /etc/fluent-bit/fluent-bit.conf
autostart=true
autorestart=true
stderr_logfile=/var/log/fluent-bit.err.log
stdout_logfile=/var/log/fluent-bit.out.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=3
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=3