# supervisor/rec_news.ini
[program:rec_news]
directory=/home/jh/code/rec-crypto-data/rcd
command=/home/jh/code/.venv/bin/python /home/jh/code/rec-crypto-data/rcd/rec_news.py
autostart=true
autorestart=true
stderr_logfile=/var/log/rec_news.err.log
stdout_logfile=/var/log/rec_news.out.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=3
stderr_logfile_maxbytes=5MB
stderr_logfile_backups=3
