# supervisor/rec_data.ini
[program:my_script]
command=/home/jh/code/.venv/bin/python /home/jh/code/rec-crypto-data/rec_python.py
autostart=true
autorestart=true
stderr_logfile=/var/log/rcd.err.log
stdout_logfile=/var/log/rcd.out.log